# fanspeed.service
[Unit]
Description=Fanspeed Daemon

[Service]
ExecStart=/usr/local/bin/fanspeed 32
Restart=on-failure

[Install]
WantedBy=multi-user.target 